# docs/launchd/com.negima.ai-fund-lab.paper-run.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.negima.ai-fund-lab.paper-run</string>

  <key>WorkingDirectory</key>
  <string>/Users/negishi/work/ai-fund-lab</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/negishi/work/ai-fund-lab/scripts/run_evening_selection.sh</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/negishi/work/ai-fund-lab/logs/paper_run.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/negishi/work/ai-fund-lab/logs/paper_run.err.log</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>16</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>16</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>16</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>16</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>16</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
  </array>
</dict>
</plist>
